# net-dialup/accel-ppp/files/pptp-init
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: $

pidfile="${PID}"
command="/usr/sbin/accel-pppd"
command_args="${ACCEL_PPPD_OPTS} -p ${PID}"



depend() {
        need net
}


start_pre() {
        checkpath -d /run/accel-ppp/

}
