# net-im/skypexmpp/files/karaka-slave
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

PIDFILE=/var/run/jabber/karaka_slave.pid

depend() {
	need net
	use jabber-server
}

start() {
	ebegin "Starting karaka-slave"
	start-stop-daemon --start --pidfile ${PIDFILE} --chuid jabber:jabber \
		--exec /usr/bin/python INSPATH/karaka-slave.py
	eend $?
}

stop() {
	ebegin "Stopping karaka-slave"
	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}
	eend $?
}
