# sys-apps/hasp/files/aksusbd.init
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need localmount
}

start() {
	ebegin "Starting HARDLOCK daemon"
	start-stop-daemon --start --quiet --exec /usr/sbin/aksusbd -- -l3
	eend $?
}

stop() {
	ebegin "Stopping HARDLOCK daemon"
	start-stop-daemon --stop --quiet --exec /usr/sbin/aksusbd
	eend $?
}
