# net-nds/ipa-server/files/ipa_webgui
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $


depend() {
	need net
	use apache
}

start() {
	
	ebegin "Starting ipa_webgui"
	 start-stop-daemon --start --quiet --pidfile /var/run/ipa_webgui.pid \
	 --exec /var/lock/subsys/ipa_webgui 

	 eend ${?}
}	

stop() {
	ebegin "Stopping ipa_webgui"	
	start-stop-daemon --stop --pidfile /var/run/ipa_webgui.pid
}	

restart() {
	stop
	start
}	

case "$1" in
  start)
  	start
	;;
  stop)
  	stop
	;;
  status)
        status $PROG
        ;;
  restart)
  	restart
	;;
  condrestart)
  	[ -f /var/lock/subsys/ipa_webgui ] && restart || :
	;;
  reload)
	exit 3
	;;
  *)
	echo $"Usage: $0 {start|stop|status|restart|condrestart}"
	exit 2
esac

exit $?
