# net-voip/sipwitch/files/sipwitch.init
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net
	use logger
}

start() {
	start-stop-daemon --start --exec /usr/sbin/sipw
}

stop() {

	start-stop-daemon --stop
}

